# restic-backup.service
# ~/.config/systemd/user/restic-backup.service
# Systemd service unit for automated restic backups to B2

[Unit]
Description=Restic Backup to Backblaze B2
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Execute the backup script
ExecStart=%h/bin/restic-backup.sh

# Restart policy for network failures
Restart=on-failure
RestartSec=5min
RestartPreventExitStatus=0

# Resource limits
Nice=19
IOSchedulingClass=idle
IOSchedulingPriority=7
CPUQuota=50%

# Logging
StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=restic-backup

# Timeout (4 hours should be enough for most backups)
TimeoutStartSec=4h

[Install]
WantedBy=default.target