# restic-backup.service
# ~/.config/systemd/user/restic-backup.service
# Systemd service unit for automated restic backups to B2

[Unit]
Description=Restic Backup to Backblaze B2
After=network-online.target
Wants=network-online.target

# Optional: Automatically send failure notification if service fails
# OnFailure=telegram-notify@%n.service

[Service]
Type=oneshot

# Execute the backup script
ExecStart=~/.config/restic/restic-backup.sh

# Restart policy for network failures
Restart=on-failure
RestartSec=5min
RestartPreventExitStatus=0

# Timeout (4 hours should be enough for most backups)
TimeoutStartSec=4h

# Optional: Automatically send success notification
# ExecStartPost=/run/media/user01/1TB-Disk/Code/Scripts/telegram-notifier/telegram-notifier %n